# recipes-connectivity/bluetooth-network/files/bluetooth-network.service
[Unit]
After=bluetooth.service
PartOf=bluetooth.service

[Service]
ExecStart=/usr/bin/br.sh

[Install]
WantedBy=bluetooth.target